# build/services.d/dbus/run
#!/usr/bin/execlineb -P
fdmove -c 2 1
dbus-daemon --system --nofork --nosyslog --nopidfile --print-address
