# build/services.d/xvfb-xmonad/run
#!/usr/bin/execlineb -P
with-contenv
importas -i nixbin NIXBIN
importas -i path PATH
importas -i userlogin USERLOGIN
cd /home/user
s6-setuidgid "${userlogin}"
s6-envuidgid -B -i "${userlogin}:${userlogin}"
backtick -n USER { echo "${userlogin}" }
backtick -n HOME { homeof "${userlogin}" }
fdmove -c 2 1
"${nixbin}/xvfb-run" --server-num=0 --server-args="-screen 0 1920x1080x24" "${nixbin}/bash" -l -c "exec '${nixbin}/xmonad'"
